# redshift-gtk.service
[Unit]
Description=Redshift, night screen filter

[Service]
ExecStart=/usr/bin/redshift-gtk

[Install]
WantedBy=default.target
